# recipes-talologger/talologger-targets/files/all-talologger-services.target
[Unit]
Description=A target that starts all home automation services

