# recipes-talologger/talologger-targets/files/talologger-visualization.target
[Unit]
Description=A target that starts home automation visualization services

