# recipes-talologger/talologger-targets/files/talologger-logging.target
[Unit]
Description=A target that starts home automation logging services

